# systemd/rog-thermal-guardian.service
[Unit]
Description=ROG Monitor — Guardian Termico CPU/GPU (fans primero, potencia despues)
Documentation=https://github.com/rog-monitor/rog-monitor
After=multi-user.target nvidia-persistenced.service
Wants=nvidia-persistenced.service

[Service]
Type=simple
# Ruta absoluta al script: ajustar si el repo vive en otra ubicacion.
# La app sustituye __ROG_MONITOR_REPO__ por la ruta real antes de copiar esta
# unidad a /etc/systemd/system/.
ExecStart=/usr/bin/env bash __ROG_MONITOR_REPO__/scripts/rog-thermal-guardian.sh
# Techo configurable sin editar esta unidad: usar un override.conf
# en /etc/systemd/system/rog-thermal-guardian.service.d/override.conf
# con Environment="ROG_THERMAL_CPU_CEILING=92" y
# Environment="ROG_THERMAL_GPU_CEILING=83". La app (main.js) escribe ese
# override.conf via pkexec cuando el usuario cambia los techos
# desde el modal de Centro de Poder.
Environment="ROG_THERMAL_CPU_CEILING=92"
Environment="ROG_THERMAL_GPU_CEILING=83"
Environment="ROG_THERMAL_INTERVAL=2"
Restart=on-failure
RestartSec=5
# Necesita escribir en sysfs (asus-armoury, platform_profile, throttle_thermal_policy)
# y leer temperatura via nvidia-smi: requiere root. NO ejecutar como usuario.
User=root
Group=root
# Endurecimiento conservador: no aislamos /sys porque el script necesita
# escribir varias rutas sysfs reales del firmware ASUS y NVIDIA.
NoNewPrivileges=true
# false: el script vive en el home del usuario (dentro del repo) y publica
# estado para la UI en ~/.local/share/rog-monitor/thermal-guardian-state.json.
# Con true/read-only systemd puede ocultar el repo o bloquear ese archivo.
ProtectHome=false
ProtectSystem=false

[Install]
WantedBy=multi-user.target
